# envsense/debian/postinst
#!/bin/bash

systemctl enable besic-envsense
systemctl start besic-envsense
