# nullsense/debian/postinst
#!/bin/bash

systemctl enable besic-nullsense
systemctl start besic-nullsense
