# dummy-csvs/debian/postinst
#!/bin/bash

systemctl enable besic-dummy-csvs
systemctl start besic-dummy-csvs
