# audio-py/debian/postinst
#!/bin/bash

systemctl enable besic-audio-py
systemctl start besic-audio-py
